# ntipa-rabbitmq/rabbitmq.conf
[program:rabbitmq]
command=/rabbitmq.sh
user=root
stopasgroup = true
redirect_stderr = true
stdout_logfile=/var/log/rabbitmq-stdout.log