# ntipa-rabbitmq/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D


[program:rabbitmq]
command=/usr/sbin/rabbitmq-server